# overlay/etc/udev/rules.d/99-gdm-hack.rules
SUBSYSTEM=="graphics", ACTION=="add", DEVPATH=="/devices/platform/display-subsystem/graphics/fb0", RUN+="/usr/lib/scripts/gdm-hack.sh"
